# root/etc/s6-overlay/s6-rc.d/ariang-conf-init/up
/etc/s6-overlay/s6-rc.d/ariang-conf-init/run