# root/etc/s6-overlay/s6-rc.d/svc-ariang/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

caddy run --config /usr/local/caddy/Caddyfile --adapter=caddyfile